# dashboard/deploy/quickinfo-dashboard-web.service
[Unit]
Description=QuickInfo Dashboard Web
After=network.target

[Service]
Type=simple
User=bitnami
WorkingDirectory=/home/bitnami/quickinfo-autoblog/dashboard/frontend
ExecStart=/usr/bin/npm start
Environment=PORT=3003
Environment=NODE_ENV=production
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
